# scripts/stellarity-instance.service
[Unit]
Description=Stellarity Instance Server
After=network.target

[Service]
Type=simple
User=stellarity
WorkingDirectory=/home/stellarity/stellarity/modules/instance
ExecStart=/home/stellarity/.bun/bin/bun run src/index.ts
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
